# .config/systemd/user/dms-sway-colors.path
[Unit]
Description=Regenerate sway colors when DMS theme changes

[Path]
PathModified=/home/alertxsto/.local/share/color-schemes/DankMatugenDark.colors
Unit=dms-sway-colors.service

[Install]
WantedBy=default.target
